# ansible/files/flask/celery_worker.conf
[program:flower_worker]
command=/srv/.virtualenv/bin/celery  -A dummy_task worker -B --loglevel=INFO
