# ansible/files/flask/web_flower_server.conf
[program:flower_web]
command=/srv/.virtualenv/bin/python  /srv/flask/flower_web.py
