# ansible/files/flask/web_flower_server.conf
[program:web_flower]
command=/srv/.virtualenv/bin/python  /srv/flask/web_flower.py
