# ansible/files/flask/web_flower_server.conf
[program:web_flower]
command=/srv/.virtualenv/bin/gunicorn -w4 -b 0.0.0.0:5000 web_flower:app
