# led/fcserver.conf
[program:fcserver]
command=/usr/local/bin/fcserver /usr/local/bin/magnus_event.json
stderr_logfile = /var/log/fcserver.stderr.log
stdout_logfile = /var/log/fcserver.stdout.log
user=root
