# led/led.conf
[program:led]
command=/bin/bash /home/pi/MagnusFlora/led/led.sh
stderr_logfile = /var/log/led.stderr.log
stdout_logfile = /var/log/led.stdout.log
user=root
