# samples/fadecandy_server.conf
[program:fadecandy]
command=/usr/local/bin/fcserver /usr/local/bin/fcserver.json
