# samples/teststrip.conf

[program:teststrip]
command=/home/pi/brian/teststrip.sh
