# rest/jarvis.conf
[program:jarvis]
command=/bin/bash /home/pi/MagnusFlora/rest/jarvis.sh
stderr_logfile = /var/log/jarvis.stderr.log
stdout_logfile = /var/log/jarvis.stdout.log
user=root
