# rest/sound.conf
[program:jarvis]
command=/bin/bash /home/pi/MagnusFlora/rest/souund.sh
stderr_logfile = /var/log/sound.stderr.log
stdout_logfile = /var/log/sound.stdout.log
user=root
