# rest/sound.conf
[program:sound]
command=/bin/bash /home/pi/MagnusFlora/rest/sound.sh
stderr_logfile = /var/log/sound.stderr.log
stdout_logfile = /var/log/sound.stdout.log
user=root
