# web/static_web.conf
[program:static_web]
command=/bin/bash /home/pi/MagnusFlora/web/static_web.sh
stderr_logfile = /var/log/mf_static_web.stderr.log
stdout_logfile = /var/log/mf_static_web.stdout.log
user=root
